# token-pool.service
[Unit]
Description=Token Pool Proxy - LLM API Token Pool Manager
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/token-pool
ExecStart=/opt/token-pool/token-pool
Environment=RUST_LOG=info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
